# rootfs/etc/services.d/twingate/finish
#!/usr/bin/with-contenv bashio

# ==============================================================================
# Cleanup Twingate Connector
# ==============================================================================

bashio::log.info "Stopping Twingate Connector service..."

# Get configuration
CONTAINER_NAME=$(bashio::config 'container_name')

# Stop and remove the Twingate container
if docker ps -q -f name="${CONTAINER_NAME}" | grep -q .; then
    bashio::log.info "Stopping Twingate container..."
    docker stop "${CONTAINER_NAME}" || true
fi

if docker ps -aq -f name="${CONTAINER_NAME}" | grep -q .; then
    bashio::log.info "Removing Twingate container..."
    docker rm "${CONTAINER_NAME}" || true
fi

bashio::log.info "Twingate Connector cleanup completed"
